# live-build/config/includes.chroot/etc/sv/NetworkManager/run
#!/bin/sh
exec /usr/sbin/NetworkManager --no-daemon
