# live-build/config/includes.chroot/etc/xdg/autostart/tenebraos-wallpaper.desktop
[Desktop Entry]
Type=Application
Name=TenebraOS Wallpaper
Comment=Set TenebraOS wallpaper
Exec=/usr/share/calamares/set-wallpaper.sh
Terminal=false
NoDisplay=true
OnlyShowIn=KDE;
StartupNotify=false
